# integrations/status-panel/agentos-status-panel.service
[Unit]
Description=AgentOS status panel — boot-health web view of the AgentOS + Nimbus stack (127.0.0.1:9123)
Documentation=file://%h/Documents/AgentOS/integrations/status-panel/README.md
After=graphical-session.target

[Service]
Type=simple
# Mostly read-only: serves panel.html + a /status.json computed from `systemctl show` and short
# port probes. Loopback-only, stdlib-only (no venv/deps). %h-relative so it's portable. The one
# write path (POST /dispatch, ADR-0039) mutates NO system state itself — it records a ledger entry
# under the runtime dir (RuntimeDirectory below) and launches an out-of-sandbox `systemd-run --user`
# worker that does the first-aid / agent investigation.
ExecStart=/usr/bin/python3 %h/Documents/AgentOS/integrations/status-panel/status_panel.py
Environment="AGENTOS_STATUS_HOST=127.0.0.1"
Environment="AGENTOS_STATUS_PORT=9123"
# Route a single debounced swaync toast on a NEW post-boot service failure (surface-labor
# contract: swaync owns the interrupt). Set to 0 to silence.
Environment="AGENTOS_STATUS_NOTIFY=1"
# Dispatch-an-agent (ADR-0039): AGENTOS_DISPATCH=0 hard-disables the whole feature;
# AGENTOS_DISPATCH_CLOUD=0 keeps local (Hermes) dispatch but removes the cloud (Claude) target,
# so no service logs ever leave the box.
Environment="AGENTOS_DISPATCH=1"
Environment="AGENTOS_DISPATCH_CLOUD=1"
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=5
KillMode=mixed
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
# The dispatch (ADR-0039) + adopt (ADR-0043) ledgers live here — RuntimeDirectory carves writable,
# 0700 dirs under $XDG_RUNTIME_DIR even though ProtectSystem=strict mounts the rest of the fs
# read-only. Preserve them across a restart so an in-flight incident/job isn't lost. The durable
# transcripts are written by the (unsandboxed) systemd-run workers under ~/.local/state.
RuntimeDirectory=agentos-dispatch agentos-adopt
RuntimeDirectoryPreserve=yes
# Hardening — it reads unit state + its own files, serves on loopback, and writes ONLY the runtime
# ledger above; the heavy lifting runs in a separate, unsandboxed systemd-run worker. (Phase-1
# on-box check: that `systemd-run --user` launches that worker from under NoNewPrivileges=true.)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true

[Install]
WantedBy=default.target
