# integrations/hermes/dashboard/hermes-dashboard.service
[Unit]
Description=Hermes Agent web dashboard (config, sessions, kanban board) on 127.0.0.1:9119
Documentation=file://%h/Documents/AgentOS/integrations/hermes/dashboard/README.md
After=network-online.target hermes-gateway.service
Wants=network-online.target
# Soft ordering only: the dashboard reads kanban.db directly and works even if the
# gateway is down, but starting after it keeps the board populated and logs tidy.

[Service]
Type=simple
WorkingDirectory=/home/corey/.hermes
# --no-open: never try to spawn a browser from the service context (no DISPLAY).
#   The UI is simply *up* at http://127.0.0.1:9119 — open it whenever you like.
# cmd_dashboard re-execs into a foreground uvicorn server and kills stale dashboard
# PIDs on start, so this is idempotent and correct for Type=simple.
ExecStart=/home/corey/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main dashboard --no-open --host 127.0.0.1 --port 9119
Environment="PATH=/home/corey/.hermes/hermes-agent/venv/bin:/home/corey/.hermes/hermes-agent/node_modules/.bin:/home/corey/.local/node-v22.22.3-linux-x64/bin:/home/corey/.local/bin:/home/corey/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/home/corey/.hermes/hermes-agent/venv"
Environment="HERMES_HOME=/home/corey/.hermes"
# First start may build the frontend (node present); give it room. Subsequent starts
# no-op the build and come up fast.
TimeoutStartSec=300
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=4
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
