# integrations/swaync-race/plasma-plasmashell.service.d/after-swaync.conf
# Nimbus/AgentOS — make swaync the notification daemon, not plasmashell.
#
# plasma-plasmashell.service is started outside swaync.service's ordering
# transaction, so swaync's own `Before=plasma-plasmashell.service` never
# constrained anything and plasmashell won the race for
# org.freedesktop.Notifications (KDE only yields the name if it's already
# taken at notification-engine init).
#
# This drop-in pulls swaync into plasmashell's start transaction and orders
# plasmashell *after* swaync's start job. With swaync's `Type=dbus`, that
# start job only completes once swaync has acquired the bus name — so
# plasmashell initializes its notification engine, finds the name taken, and
# politely defers.
#
# `Wants=` (not `Requires=`) keeps this fail-open: if swaync fails to start,
# plasmashell still comes up and reclaims notifications — the desktop is never
# left without a notification daemon.
#
# Reversible: `rm` this file + `systemctl --user daemon-reload` restores stock
# Plasma notification ownership.
[Unit]
Wants=swaync.service
After=swaync.service
