# integrations/swaync-race/swaync.service.d/nimbus-race.conf
# Nimbus/AgentOS — let swaync reliably win org.freedesktop.Notifications on boot.
#
# Root cause v3 (2026-06-23, from a live boot journal + a busctl-monitor trace of
# swaync's own D-Bus traffic). swaync is a GTK4 + libadwaita app. During init,
# BEFORE it requests the org.freedesktop.Notifications bus name, libadwaita issues
# SYNCHRONOUS org.freedesktop.portal.Settings Read/ReadAll calls (color-scheme,
# accent-color, …). On a cold boot xdg-desktop-portal is not yet up — in fact this
# box's xdg-desktop-portal.service itself fails its first start (its backend isn't
# ready inside its 15s budget) and isn't back for ~3 minutes — so each of those
# synchronous portal calls blocks the full ~25s D-Bus reply timeout. A monitor
# trace showed THREE such libadwaita Read calls ahead of RequestName: ~75s of
# stall, which is exactly why v1's 45s TimeoutStartSec "was necessary but not
# sufficient." swaync never acquires the name in time, its start job fails, and a
# *failed* start job still satisfies plasmashell's `After=swaync` ordering — so
# plasmashell comes up, grabs the name, and keeps it for the session (swaync
# requests the name without replacement and can't reclaim it).
#
# (v2, SUPERSEDED and REMOVED: an ExecStartPre that *waited for the portal to
# answer* before starting swaync. It made swaync's success depend on the one
# service guaranteed to be dead at cold boot — its busy-wait worst case (60 ×
# (2s busctl timeout + 0.5s) = 150s) overran TimeoutStartSec=90, so start-pre was
# KILLED and swaync's daemon never even ran. Strictly counterproductive: it
# *guaranteed* the failure it was meant to prevent, every boot.)
#
# Fix v3 — remove the portal from swaync's path instead of waiting on it:
#   1. Environment=ADW_DISABLE_PORTAL=1 — tell libadwaita not to query the Settings
#      portal at all. This deletes the blocking color-scheme/accent reads (trace:
#      10 portal calls → 1, and the survivor is GTK4's own async query, which does
#      NOT gate name acquisition). swaync reaches RequestName promptly even with a
#      dead portal, so it wins the name AND doesn't stall plasmashell (which waits
#      behind it via After=swaync). swaync's look is driven by its own CSS, not the
#      portal color-scheme, so nothing visible changes.
#   2. The Wayland-socket ExecStartPre stays (swaync is a layer-shell client; it
#      must not launch before its compositor socket exists). Bounded + fail-open.
#   3. No Wants=/After= on xdg-desktop-portal — swaync no longer depends on it.
#   4. TimeoutStartSec=45 — generous margin for cold GTK init; far more than the
#      <3s warm acquire now that the portal stall is gone. If swaync still somehow
#      fails to acquire within the budget, plasmashell reclaims notifications
#      (fail-open — the desktop is never left without a notification daemon).
#
# Reversible: integrations/swaync-race/restore.sh (or rm this file + daemon-reload).
# The underlying xdg-desktop-portal cold-boot failure is a separate disease (it also
# breaks plasmashell portal registration + screencast at boot); removing this
# cascade should let it heal, but it is worth its own follow-up.

[Unit]
StartLimitIntervalSec=120
StartLimitBurst=8

[Service]
Environment=ADW_DISABLE_PORTAL=1
TimeoutStartSec=45
RestartSec=2
# Don't launch until the Wayland socket exists (swaync is a layer-shell client).
# Bounded ~40s busy-wait, then exit 0 REGARDLESS (fail-open).
ExecStartPre=/usr/bin/bash -c 'sock="${XDG_RUNTIME_DIR}/${WAYLAND_DISPLAY}"; for i in $(seq 1 80); do [ -S "$sock" ] && exit 0; sleep 0.5; done; exit 0'
