# integrations/swaync-race/swaync.service.d/nimbus-race.conf
# Nimbus/AgentOS — let swaync reliably win org.freedesktop.Notifications on boot.
#
# Root cause v2 (2026-06-22, from a live boot journal). swaync queries
# org.freedesktop.portal.Settings during init, BEFORE it acquires the bus name.
# On a cold boot the xdg-desktop-portal frontend is D-Bus-activated and not yet
# warm, so swaync's own Settings call blocks ~25s on the portal timeout; with GTK
# init on top, the total start blows past TimeoutStartSec and systemd ABORTS the
# start job. A *failed* start job still satisfies plasmashell's `After=swaync`
# ordering, so plasmashell comes up and grabs org.freedesktop.Notifications while
# it's free; every later swaync restart then finds the name taken and start-limit-
# bricks for the session (plasmashell keeps notifications — fail-open).
#
# (v1, superseded but kept: a short DefaultTimeoutStartUSec=15s aborted the cold
# start; bumping to 45s + a Wayland-socket wait was necessary but NOT sufficient,
# because the dominant stall is the portal, not the socket.)
#
# Fix — take the portal off swaync's critical path so its FIRST start acquires the
# name promptly; then plasmashell (After=swaync, swaync Type=dbus) finds the name
# taken and defers, deterministically:
#   1. Wants=xdg-desktop-portal.service — pull the portal frontend into the
#      transaction so it warms alongside swaync. We deliberately do NOT
#      `After=xdg-desktop-portal.service` (v3, 2026-06-22): the portal is itself
#      ordered `After=graphical-session.target`, and stock `plasma-core.target` is
#      `After=plasma-plasmashell.service`, so `plasmashell → After swaync → After
#      portal → … → plasma-core.target → plasmashell` formed an ordering CYCLE.
#      systemd breaks a cycle by DELETING a job — it deleted plasmashell's start
#      job, so the entire shell (panels + dock) never launched at boot. The
#      ExecStartPre busy-wait in (2) already sequences swaync after the portal
#      *answers* at runtime — more robust than systemd ordering and cycle-free —
#      so Wants= (a pure dependency, no ordering edge) is all we need.
#   2. ExecStartPre (after the Wayland-socket wait) blocks until the portal Settings
#      interface actually ANSWERS, so swaync's own query returns instantly instead
#      of hitting its ~25s internal timeout. Bounded ~30s, then `exit 0` REGARDLESS
#      (fail-open: a dead portal can never brick swaync).
#   3. TimeoutStartSec=90 — margin for the one-time cold warmup (now paid in
#      ExecStartPre) plus swaync's measured <2s warm acquire. Warm boots pay ~0.
#
# Reversible: integrations/swaync-race/restore.sh (or rm this file + daemon-reload).
# Fail-open: if swaync still never starts, plasmashell reclaims notifications.

[Unit]
# Wants= (NOT After=) — pulling the portal in is fine, but ordering swaync after it
# closes a boot cycle that deletes plasmashell's start job. See note (1) above.
Wants=xdg-desktop-portal.service
StartLimitIntervalSec=120
StartLimitBurst=8

[Service]
TimeoutStartSec=90
RestartSec=2
# (1) don't launch until the Wayland socket exists (swaync is a layer-shell client).
ExecStartPre=/usr/bin/bash -c 'sock="${XDG_RUNTIME_DIR}/${WAYLAND_DISPLAY}"; for i in $(seq 1 80); do [ -S "$sock" ] && exit 0; sleep 0.5; done; exit 0'
# (2) warm + wait for the xdg-desktop-portal Settings interface to actually answer,
#     so swaync's own Settings query doesn't hit its ~25s internal timeout and blow
#     the start budget. Bounded (~30s) + fail-open; short per-call busctl timeout.
ExecStartPre=/usr/bin/bash -c 'for i in $(seq 1 60); do /usr/bin/busctl --user --timeout=2 call org.freedesktop.portal.Desktop /org/freedesktop/portal/desktop org.freedesktop.portal.Settings ReadAll as 0 >/dev/null 2>&1 && exit 0; sleep 0.5; done; exit 0'
